# example.conf.d/sleeper.ini
[program:sleeper]
command = bash -c 'sleep 4s'
autorestart = true
redirect_stderr = true
