# example.conf.d/truther.ini
[program:truther]
command = false
autorestart = unexpected
redirect_stderr = true
