# tests/one-failure.ini
[program:one-failure]
command = ./bin/procwatch-tester -t 2s -s 1

