# tests/one-success.ini
[program:one-success]
command = ./bin/procwatch-tester -t 2s
autorestart = false