# tests/one-success-cron.ini
[procwatch]
;stdout_logfile = stdout

[program:one-success-cron]
command = ./bin/procwatch-tester -t 2s
schedule = "*/30 * * * * *"