# tests/one-success-nondefault-exit-code.ini
[program:one-success-exitcode]
command = ./bin/procwatch-tester -t 2s -s 42
exitcodes = 0,2,42

