# tests/one-success.ini
[procwatch]
logfile = /tmp/procwatch.log

[program:one-success]
command = cat /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab /etc/fstab
autorestart = true
startretries = 10000000
stdout_logfile_maxbytes = 1MB
stdout_logfile_backups = 1