# tests/one-failure-ar.ini
[program:one-failure-ar]
command = ./bin/procwatch-tester -t 50ms -s 1
autorestart = true
startretries = 3
