# tests/one-success.ini
[program:one]
command = ./bin/procwatch-tester

