# tests/one-success-cron.ini
[program:one-success-cron]
command = ./bin/procwatch-tester -t 2s
schedule = "* * * * *"